# files/sup-nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
stdout_events_enabled=true
stderr_events_enabled=true

[program:testssl.sh-webfrontend]
command=/usr/local/bin/testssl.sh-webfrontend/SSLTestPortal.py
